# claude自動引当ツール/deploy/systemd/user/zaiko-yahoo-orders.service
[Unit]
Description=Zaiko Yahoo orders import
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/ubuntu/zaiko-tool/app
EnvironmentFile=/home/ubuntu/zaiko-tool/app/.env
ExecStart=/home/ubuntu/zaiko-tool/app/scripts/job_yahoo_orders.sh
